# Cleanup/cleanup.crontab
10 10,22 * * * . ~/.bash_profile ; python ndeCloudCleanup.py
1 1,3,5,7,9,11,13,15,17,19,21,23 * * * . ~/.bash_profile ; python ndeCloudDBCleanup.py >> dbCleanup.log
40 * * * * . ~/.bash_profile ; python cleanupElasticSearch.py >> esCleanup.log

